# twitter-bot.service
[Unit]
Description=Twitter URL to Image Telegram Bot
Documentation=https://github.com/harusClaw/twitter-capture
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/openclaw/.nanobot/workspace/projects/twitter-capture

# Environment setup
Environment="PATH=/home/openclaw/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=-/home/openclaw/.nanobot/workspace/projects/twitter-capture/.env

# Execution - use module entry point
ExecStart=/home/openclaw/.local/bin/uv run python -m twitter_capture
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=twitter-capture-bot

[Install]
WantedBy=default.target
